# debian/prerm
#!/bin/bash

systemctl disable maratona-firstboot.service
